# deploy/systemd/prometheus-node-exporter.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/prometheus-node-exporter \
  --web.listen-address=127.0.0.1:9100 \
  --collector.systemd
