# deploy/systemd/blocknet-pool-alertmanager-discord-relay.service
[Unit]
Description=Blocknet Pool Alertmanager Discord Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blocknet
Group=blocknet
WorkingDirectory=/opt/blocknet/blocknet-pool
EnvironmentFile=-/etc/blocknet/monitoring.env
ExecStart=/usr/bin/python3 /opt/blocknet/blocknet-pool/scripts/alertmanager_discord_relay.py --listen 127.0.0.1:24810
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadOnlyPaths=/opt/blocknet/blocknet-pool

[Install]
WantedBy=multi-user.target
