# deploy/systemd/blocknetd.service
[Unit]
Description=Blocknet Core Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blocknet
Group=blocknet
WorkingDirectory=/opt/blocknet/blocknet-core/current
Environment=BLOCKNET_P2P_KEY=/var/lib/blocknet/data/identity.key
ExecStart=/opt/blocknet/blocknet-core/current/blocknet --daemon --api 127.0.0.1:8332 --wallet /var/lib/blocknet/wallet.dat --data /var/lib/blocknet/data --listen /ip4/0.0.0.0/tcp/28080
Restart=always
RestartSec=5
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/blocknet/blocknet-core /var/lib/blocknet

[Install]
WantedBy=multi-user.target
