# deploy/systemd/postgresql@.service.d/restart-blocknet-pool.conf
[Service]
# Pool services keep local state derived from live Postgres connections. Bounce
# them after the cluster comes back so monitors and workers reconnect cleanly
# after package upgrades or manual database restarts.
ExecStartPost=/usr/bin/systemctl try-restart blocknet-pool-api.service blocknet-pool-stratum.service blocknet-pool-monitor.service
