# deploy/systemd/postgresql@.service.d/restart-blocknet-pool.conf
[Service]
# The pool's current Postgres wrapper can wedge on externally restarted
# connections. When a PostgreSQL cluster starts again, bounce the pool so it
# reconnects cleanly after package upgrades or manual database restarts.
ExecStartPost=/usr/bin/systemctl try-restart blocknet-pool-api.service blocknet-pool-stratum.service
